# devel/comum/crontabs
20 23,08 * * * /usr/local/scripts/atualiza.sh
